# internet-stop.service
[Unit]
Description=Block internet for listed hosts

[Service]
Type=forking
ExecStart=/home/pi/internet-block/arp-spoof-start
