# internet-start.service
[Unit]
Description=Block internet for listed hosts

[Service]
Type=oneshot
ExecStart=/home/pi/internet-block/arp-spoof-start
